# Content/DefenderCustomDetections/Endpoint/PersistenceViaScheduledTask.yaml
displayName: Persistence via scheduled task creation with suspicious command
isEnabled: true
queryCondition:
  queryText: |
    DeviceProcessEvents
    | where Timestamp > ago(1h)
    | where FileName == "schtasks.exe"
    | where ProcessCommandLine has "/create"
    | where ProcessCommandLine has_any ("powershell", "cmd.exe", "wscript",
        "cscript", "mshta", "rundll32", "regsvr32", "certutil",
        "bitsadmin", "msiexec")
    | where ProcessCommandLine !has "Microsoft\\Windows\\Defrag"
    | where ProcessCommandLine !has "GoogleUpdate"
    | project Timestamp, DeviceId, DeviceName, ReportId,
              FileName, ProcessCommandLine, FolderPath,
              InitiatingProcessFileName, InitiatingProcessCommandLine,
              InitiatingProcessSHA1,
              InitiatingProcessAccountUpn, InitiatingProcessAccountSid
  lastModifiedDateTime: "2026-03-23T00:00:00Z"
schedule:
  period: "1H"
detectionAction:
  alertTemplate:
    title: "Persistence via scheduled task creation with suspicious command"
    description: |
      A scheduled task was created that executes a script interpreter or
      living-off-the-land binary. Attackers create scheduled tasks to maintain
      persistence and execute payloads on a recurring schedule.
    severity: medium
    category: Persistence
    mitreTechniques:
      - T1053.005
    recommendedActions: |
      1. Review the scheduled task configuration and the command it executes.
      2. Check the parent process that created the task for legitimacy.
      3. Verify whether the task was created by an administrator or automated tool.
      4. If suspicious, delete the scheduled task and scan the device.
      5. Check for other persistence mechanisms on the same device.
    impactedAssets:
      - "@odata.type": "#microsoft.graph.security.impactedDeviceAsset"
        identifier: deviceId
  responseActions:
    - "@odata.type": "#microsoft.graph.security.collectInvestigationPackageResponseAction"
      identifier: deviceId
    - "@odata.type": "#microsoft.graph.security.runAntivirusScanResponseAction"
      identifier: deviceId
